# xrandr/99-monitor-connect.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/etc/xrandr.d/resolve-screen"
