# systemd/org-sync/org-sync.slice
# -*- mode: conf -*-

[Unit]
Description=Limited resources Slice
DefaultDependencies=no
Before=slices.target

[Slice]
CPUQuota=50%
MemoryLimit=1G
